# recipes-extended/psplash/files/psplash-start-override.service
[Unit]
Description=Start psplash boot splash screen
DefaultDependencies=no
RequiresMountsFor=/run
ConditionPathExists=!/run/psplash_fifo
